# passhub-backup-scripts/etc/cron.daily/passhub-backup
#!/bin/sh

/var/lib/passhub.backup/do_backup.sh

